# scripts/btkeepalive.conf
start on started btm
stop on stopping btm

script
    exec /bin/sh /mnt/us/btconnect.sh
end script
